# systemd/trading_bot_v9_oneshot.timer
[Unit]
Description=Run V9 every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
