# systemd/trading_bot_v9_oneshot.service
[Unit]
Description=V9_REGIME_CORE One-shot Runner
After=network-online.target

[Service]
Type=oneshot
User=trader
WorkingDirectory=/home/trader/trading_system
EnvironmentFile=-/home/trader/trading_system/.env
Environment=V9_LIVE_MODE=LIVE
ExecStart=/home/trader/ops/run_v9_hourly.sh
StandardOutput=journal
StandardError=journal
