# trading_system/trading_bot.service
[Unit]
Description=V8 A-only Trading Bot (1D macro + vol regime sizing)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/trading_system
ExecStart=/usr/bin/python3 -u -m bots.bot_a.main_bot_a
Restart=no
RestartSec=30
StandardOutput=journal
StandardError=journal

# V8: C engine deprecated. No auto-restart per governance.
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
