# trading_system/trading_bot.service
[Unit]
Description=BNB/USDT Futures Testnet (futures_run.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/trading_system
ExecStart=/usr/bin/python3 -u -m bots.bot_c.futures_run
Restart=always
RestartSec=30
StandardOutput=append:/home/trader/trading_system/logs/paper_out.log
StandardError=append:/home/trader/trading_system/logs/paper_err.log

# 避免記憶體不足時頻繁重啟
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
